# deploy/display-on.service
# /etc/systemd/system/display-on.service
# Turns on the monitor in the morning. Paired with display-on.timer.

[Unit]
Description=Turn on display

[Service]
Type=oneshot
User=pi
Environment=DISPLAY=:0
ExecStart=/usr/bin/xset dpms force on
